# scripts/com.centralhardware.battery-health.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.centralhardware.battery-health</string>

    <key>ProgramArguments</key>
    <array>
        <string>/PATH/TO/YOUR/batteryStatServer/scripts/report_battery_health.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>BATTERY_SERVER_URL</key>
        <string>http://localhost:8080</string>
    </dict>

    <!-- Run every hour (3600 seconds) -->
    <key>StartInterval</key>
    <integer>3600</integer>

    <!-- Run at load -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Logs -->
    <key>StandardOutPath</key>
    <string>/tmp/battery_health.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/battery_health_error.log</string>
</dict>
</plist>
